# scripts/com.cel.docs.news.fetch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cel.docs.news.fetch</string>

  <key>WorkingDirectory</key>
  <string>/Applications/XAMPP/xamppfiles/htdocs/docs</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>npm run news:fetch</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Applications/XAMPP/xamppfiles/htdocs/docs/logs/news-fetch.log</string>
  <key>StandardErrorPath</key>
  <string>/Applications/XAMPP/xamppfiles/htdocs/docs/logs/news-fetch.err.log</string>

  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
